# ops/launchd/bank_rates_scraper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>se.banks.listed.rates.scraper</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Volumes/ExternalSSD/dataengineering/swedish-mortgages-dbt/ingestion/bank_rates_scraper.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/swedish_mortgages_bank_rates_scraper.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/swedish_mortgages_bank_rates_scraper.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>

